# RConfig/server1/rc.local
###e
echo "Server2"
echo "Update by omid"
###
